# data/falcond.service
[Unit]
Description=Falcon Daemon for Game Optimization
After=network.target dbus.service
Requires=dbus.service

[Service]
Type=simple
ExecStart=/usr/bin/falcond
Restart=always
RestartSec=5

# Managed directories
RuntimeDirectory=falcond
RuntimeDirectoryMode=0755
StateDirectory=falcond
StateDirectoryMode=0755
LogsDirectory=falcond
LogsDirectoryMode=0755

# Resource limits
LimitNOFILE=65536
LimitNPROC=65536

# Security hardening (running as root but restricted)
PrivateTmp=yes
ProtectSystem=full
ProtectHome=read-only

[Install]
WantedBy=multi-user.target
